# infinite-live.service
[Unit]
Description=Infinite Game Live Server (WebSocket)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/Infinite-Game-Research
Environment=OMP_NUM_THREADS=1
Environment=MKL_NUM_THREADS=1
Environment=OPENBLAS_NUM_THREADS=1
Environment=NUMEXPR_NUM_THREADS=1
Environment=VECLIB_MAXIMUM_THREADS=1
Environment=IG_SEED=0
Environment=IG_ADJUST_INTERVAL=2000
Environment=IG_LOG_EVERY=20000
Environment=IG_PRINT_EVERY_SEC=10
Environment=IG_TRAJ_WIN=20000
Environment=IG_CPX_WIN=20000
Environment=IG_DOWNSAMPLE=10
Environment=PORT=8000
ExecStart=/root/Infinite-Game-Research/.venv/bin/python /root/Infinite-Game-Research/experiments/live/server.py
Restart=always
RestartSec=3
# 日志输出到文件（避免journald限制）
StandardOutput=append:/root/Infinite-Game-Research/logs/infinite-live.out.log
StandardError=append:/root/Infinite-Game-Research/logs/infinite-live.err.log
# 如果使用journald，可以添加日志上限（可选）
# StandardOutput=journal
# StandardError=journal
# LogRateLimitIntervalSec=60
# LogRateLimitBurst=1000

[Install]
WantedBy=multi-user.target
